# infra/deployment/systemd/chad-queue-worker.service
[Unit]
Description=Chad-Core Queue Worker
Documentation=https://github.com/chad-core/chad-core
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=simple
User=chad
Group=chad
WorkingDirectory=/opt/chad-core

# Set environment
Environment="PYTHONUNBUFFERED=1"
Environment="REDIS_URL=redis://localhost:6379/0"
Environment="DATABASE_URL=postgresql+asyncpg://postgres:postgres@localhost:5432/chad_core"

# Override environment from file (optional)
EnvironmentFile=-/etc/chad-core/queue-worker.env

# Execute worker
ExecStart=/opt/chad-core/venv/bin/python -m apps.queue_worker.main

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
CPUQuota=200%
MemoryLimit=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-queue-worker

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chad-core/logs

[Install]
WantedBy=multi-user.target
